# reserveRecord.xml
<ns0:Task xmlns:ns0="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
  <ns0:RegistrationInfo>
    <ns0:Date>2024-08-20T14:23:39.6303956</ns0:Date>
    <ns0:Author>AKIOKA-LTD\Akioka-Shared</ns0:Author>
    <ns0:URI>\WebMeeting\test082002</ns0:URI>
  </ns0:RegistrationInfo>
  <ns0:Triggers>
    <ns0:TimeTrigger>
      <ns0:StartBoundary>2024-08-21T11:05:00</ns0:StartBoundary>
      <ns0:EndBoundary>2025-08-21T11:05:00</ns0:EndBoundary>
      <ns0:Enabled>true</ns0:Enabled>
    </ns0:TimeTrigger>
  </ns0:Triggers>
  <ns0:Principals>
    <ns0:Principal id="Author">
      <ns0:UserId>S-1-5-21-1051331784-2990698068-3355746272-1186</ns0:UserId>
      <ns0:LogonType>InteractiveToken</ns0:LogonType>
      <ns0:RunLevel>HighestAvailable</ns0:RunLevel>
    </ns0:Principal>
  </ns0:Principals>
  <ns0:Settings>
    <ns0:MultipleInstancesPolicy>IgnoreNew</ns0:MultipleInstancesPolicy>
    <ns0:DisallowStartIfOnBatteries>true</ns0:DisallowStartIfOnBatteries>
    <ns0:StopIfGoingOnBatteries>true</ns0:StopIfGoingOnBatteries>
    <ns0:AllowHardTerminate>true</ns0:AllowHardTerminate>
    <ns0:StartWhenAvailable>false</ns0:StartWhenAvailable>
    <ns0:RunOnlyIfNetworkAvailable>false</ns0:RunOnlyIfNetworkAvailable>
    <ns0:IdleSettings>
      <ns0:StopOnIdleEnd>true</ns0:StopOnIdleEnd>
      <ns0:RestartOnIdle>false</ns0:RestartOnIdle>
    </ns0:IdleSettings>
    <ns0:AllowStartOnDemand>true</ns0:AllowStartOnDemand>
    <ns0:Enabled>true</ns0:Enabled>
    <ns0:Hidden>false</ns0:Hidden>
    <ns0:RunOnlyIfIdle>false</ns0:RunOnlyIfIdle>
    <ns0:WakeToRun>false</ns0:WakeToRun>
    <ns0:ExecutionTimeLimit>PT72H</ns0:ExecutionTimeLimit>
    <ns0:DeleteExpiredTaskAfter>PT0S</ns0:DeleteExpiredTaskAfter>
    <ns0:Priority>7</ns0:Priority>
  </ns0:Settings>
  <ns0:Actions Context="Author">
    <ns0:Exec>
      <ns0:Command>C:\Users\Akioka-Shared\AppData\Local\Programs\Python\Python312\python.exe</ns0:Command>
      <ns0:Arguments>StartOrFinishRecord.py</ns0:Arguments>
      <ns0:WorkingDirectory>C:\python\python-auto-meeting</ns0:WorkingDirectory>
    </ns0:Exec>
  </ns0:Actions>
</ns0:Task>